# etc/profile.d/xdg-local.sh
export XDG_DATA_DIRS="$XDG_DATA_DIRS:$HOME/.local/share"
